# camera-capture.service
[Unit]
Description=Camera Capture System with GPIO Triggers
After=network.target network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/home/pi/camera_capture_system
ExecStart=/usr/bin/python3 /home/pi/camera_capture_system/main.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Environment variables (optional - you can also use .env file)
# Environment="GPIO_ENABLED=true"
# Environment="UPLOAD_ENABLED=true"

# Increase file descriptor limits if needed
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target

